# templates/supervisord.conf
[program:envd]
command=envd -c path/to/configs -t accesstoken
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/envd.log
stderr_logfile=/var/log/envd.log
user=root